# containers/laravel/mains/7.0/config/supervisor/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stopwaitsecs=3600
